# usr/lib/systemd/system/restic-check@.timer
[Unit]
Description=Check restic backup for errors on a schedule

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
