# Library/LaunchAgents/com.github.erikw.restic-automatic-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- See launchd.plist(5)for documentation on this file. -->
<!-- See https://www.launchd.info/ for a tutorial. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>com.github.erikw.restic-automatic-backup</string>
		<key>EnvironmentVariables</key>
		<dict>
			<key>PATH</key>
			<string>/usr/local/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		</dict>
		<key>ProgramArguments</key>
		<array>
			<string>/bin/bash</string>
			<string>-c</string>
			<string>export PATH=/usr/local/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin &amp;&amp; source /usr/local/etc/restic/erikw.env.sh &amp;&amp; PATH=/usr/local/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin /usr/local/bin/restic_backup.sh</string>
		</array>
		<!-- how use $HOME? -->
		<key>StandardOutPath</key>
		<string>/Users/erikw/Library/Logs/restic/restic_stdout.log</string>
		<key>StandardErrorPath</key>
		<string>/Users/erikw/Library/Logs/restic/restic_stderr.log</string>
		<key>RunAtLoad</key>
		<true/>
		<!-- Will schedule backup every day at 19:00 -->
		<key>StartCalendarInterval</key>
		<array>
			<dict>
				<key>Hour</key>
				<integer>19</integer>
				<key>Minute</key>
				<integer>0</integer>
			</dict>
		</array>
	</dict>
</plist>
