# Library/LaunchAgents/com.github.erikw.restic-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>com.github.erikw.restic-check</string>
		<key>ProgramArguments</key>
		<array>
			<string>/bin/bash</string>
			<string>-c</string>
			<string>source {{ INSTALL_PREFIX }}/etc/restic/default.env.sh &amp;&amp; {{ INSTALL_PREFIX }}/bin/restic_check.sh >>$HOME/$LOG_OUT 2>>$HOME/$LOG_ERR</string>
		</array>
		<key>EnvironmentVariables</key>
		<dict>
			<key>PATH</key>
			<string>{{ INSTALL_PREFIX }}/bin:/usr/local/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
			<key>LOG_OUT</key>
			<string>/Library/Logs/restic/check_stdout.log</string>
			<key>LOG_ERR</key>
			<string>/Library/Logs/restic/check_stderr.log</string>
		</dict>
		<key>RunAtLoad</key>
		<true/>
		<!-- Will schedule check every 1st of the month at 20:00 -->
		<key>StartCalendarInterval</key>
		<array>
			<dict>
				<key>Day</key>
				<integer>1</integer>
				<key>Hour</key>
				<integer>20</integer>
				<key>Minute</key>
				<integer>0</integer>
			</dict>
		</array>
	</dict>
</plist>
